# BillServer/scripts/bill.service
[Unit]
Description=bill service

[Service]
User=root
#Type=forking
ExecStart=/opt/bill/start.sh

